# scrapyd-supervisor.conf
[supervisord]
nodaemon=true

[program:scrapyd]
command=sh -c "scrapyd"

[program:scrapyd-deploy]
command=sh -c "scrapyd-deploy default"

[program:scrapydweb]
command=sh -c "scrapydweb"
